# systemd/ost-web.service
[Unit]
Description=Open Science Tracker Web Application
After=network.target postgresql.service

[Service]
Type=exec
User=xeradb
Group=xeradb
WorkingDirectory=/var/www/ost
Environment=DJANGO_SETTINGS_MODULE=ost_web.settings
ExecStart=/var/www/ost/ost_env/bin/gunicorn ost_web.wsgi:application --bind 127.0.0.1:8000 --workers 3 --timeout 120
ExecReload=/bin/kill -HUP $MAINPID
Restart=always
RestartSec=5

# Logging
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target 